# challs_finished/Odayaka Waters/odayaka_waters/config/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock
chmod=0770
chown=root:editor

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[supervisord]
user=root
nodaemon=true
stderr_logfile=/dev/null
stdout_logfile=/dev/null
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

[program:artisan]
command=php artisan serve --host=127.0.0.1 --port=8000
user=challenger
directory=/app/main_app
autostart=true
autorestart=true
stdout_logfile=/var/log/php/access.log
stderr_logfile=/var/log/php/error.log
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0
stopasgroup=true
killasgroup=true


[program:caddy]
command=caddy run --config /etc/caddy/Caddyfile --resume
user=root
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0
stderr_logfile_maxbytes=0

;[program:editor-app]
;command=python3 app.py
;user=editor
;directory=/app/editor/backend
;autostart=true
;autorestart=true
;stdout_logfile=/app/editor/backend/1.out
;stderr_logfile=/app/editor/backend/1.err
;stdout_logfile_maxbytes=0
;stderr_logfile_maxbytes=0
